# composer/.config/zsh/profile.d/composer.zsh
#!/usr/bin/env zsh

[ -d "$COMPOSER_HOME" ] && export PATH="$COMPOSER_HOME/vendor/bin:$PATH"
